# server-info/systemd/photoprism-ledger.timer.d/description.conf.template
[Unit]
Description=Weekly PhotoPrism Ledger Pipeline
